# dynamic_cpu_limiter.service
[Unit]
Description=Dynamic CPU Limiter Service
After=network.target

[Service]
ExecStart=/usr/local/bin/dynamic_cpu_limiter.sh
Restart=always
RestartSec=5s
User=root

[Install]
WantedBy=multi-user.target
